# nolongerbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nolongerbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/flippackstation5/python_scripts/nolonger_bot/venv/bin/python3</string>
        <string>/Users/flippackstation5/python_scripts/nolonger_bot/find_orders.py</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardErrorPath</key>
    <string>/Users/flippackstation5/python_scripts/nolonger_bot/logs/find_orders.err</string>
    <key>StandardOutPath</key>
    <string>/Users/flippackstation5/python_scripts/nolonger_bot/logs/find_orders.out</string>
</dict>
</plist>